# crontab
* * * * * /Users/kirill/bin/autocommit /Users/kirill/Personal/vimwiki/
