# crontab
# Don't forget to reload with
# crontab ~/.dotfiles/crontab

* * * * * /Users/kirill/bin/autocommit /Users/kirill/Personal/vimwiki/
